# kali-agent/kali-agent.service
[Unit]
Description=Kali Agent - AI-powered security assistant
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/kali-agent
EnvironmentFile=/opt/kali-agent/.env
ExecStart=/opt/kali-agent/.venv/bin/python /opt/kali-agent/daemon.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
